# users/sam/runit/sv/ssh-agent/run
#!/usr/bin/bash

exec ssh-agent -Da "~/.ssh/ssh-agent.sock" > /dev/null
